# debian/postinst
#!/bin/bash
set -e

FROXLOR_DIR="/var/www/froxlor"
ENV_FILE="$FROXLOR_DIR/.env"
GREEN='\033[0;32m'
YELLOW='\033[1;33m'
CYAN='\033[0;36m'
NC='\033[0m'

case "$1" in
    configure)

        # --- Dateiberechtigungen setzen ---
        chown -R www-data:www-data "$FROXLOR_DIR"
        chmod -R 755 "$FROXLOR_DIR"
        chmod -R 775 "$FROXLOR_DIR/storage"
        chmod -R 775 "$FROXLOR_DIR/bootstrap/cache"

        # --- .env erstellen wenn nicht vorhanden ---
        if [ ! -f "$ENV_FILE" ]; then
            cp "$FROXLOR_DIR/.env.example" "$ENV_FILE"
            chown www-data:www-data "$ENV_FILE"
            chmod 640 "$ENV_FILE"
        fi

        # --- App-Key generieren wenn nicht gesetzt ---
        if ! grep -q "^APP_KEY=base64:" "$ENV_FILE" 2>/dev/null; then
            sudo -u www-data php "$FROXLOR_DIR/artisan" key:generate --force --quiet 2>/dev/null || true
        fi

        # --- Systemd Queue Worker ---
        if command -v systemctl >/dev/null 2>&1; then
            systemctl daemon-reload
            systemctl enable froxlor-queue.service 2>/dev/null || true
        fi

        echo ""
        echo -e "${GREEN}╔══════════════════════════════════════════════════╗${NC}"
        echo -e "${GREEN}║       Froxlor wurde erfolgreich installiert!     ║${NC}"
        echo -e "${GREEN}╚══════════════════════════════════════════════════╝${NC}"
        echo ""
        echo -e "${CYAN}Nächste Schritte:${NC}"
        echo ""
        echo -e "  ${YELLOW}1.${NC} Datenbank konfigurieren:"
        echo "     nano /var/www/froxlor/.env"
        echo ""
        echo -e "  ${YELLOW}2.${NC} Froxlor einrichten (Migration + Admin erstellen):"
        echo "     froxlor-setup"
        echo ""
        echo -e "  ${YELLOW}3.${NC} Webserver konfigurieren:"
        echo "     Apache: /etc/froxlor/apache.conf.example"
        echo "     Nginx:  /etc/froxlor/nginx.conf.example"
        echo ""
        echo -e "  ${YELLOW}4.${NC} Panel aufrufen: http://IHR-SERVER/admin"
        echo ""
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
    ;;
esac

exit 0
